# crontab
*/15 * * * * python /app/archive_albums.py
